# persistence/home/kali/.config/autostart/start.desktop
[Desktop Entry]
Name=startup
Exec=/home/kali/scripts/startup.sh
Terminal=false
Type=Application
