# root/.config/autostart/start.desktop
[Desktop Entry]
Name=startup
Exec=/root/scripts/startup.sh
Terminal=false
Type=Application
